# systemd/updatebtw-update.service
[Unit]
Description=updatebtw automatic system update
Documentation=https://github.com/anomalyco/updatebtw
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/updatebtw update
Environment=SYSTEMD_INHIBIT=1
Environment=UPDATERBTW_ROOT=/usr/lib/updatebtw
StandardOutput=journal
StandardError=journal
